# k.yyup.com/client/wsl-network.service
[Unit]
Description=WSL Network Initialization Service
After=network.target
Wants=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/init-wsl-network.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 